# config/root/etc/profile.d/drucker.sh
#!/bin/bash
#lpstat -s -e -d -a 
lpoptions -d EPSON_XP_342_343_345_Series >/dev/null

